# units/funk_dynamic.service.in
[Unit]
Description=Funky daemon
Requires=funk.socket

[Service]
Type=notify
ExecStart=@BINDIR@/funk_dynamic

[Install]
WantedBy=multi-user.target
Also=funk.socket
