# units/funk_ro.socket.in
[Socket]
ListenStream=@RUNDIR@/funk/socket_ro
SocketMode=0400
Service=funk.service

[Install]
WantedBy=sockets.target
