# units/funk_rw.socket.in
[Socket]
ListenStream=@RUNDIR@/funk/socket
SocketMode=0600
Service=funk.service

[Install]
WantedBy=sockets.target
